# aishield/aishield.service
[Unit]
Description=AIShield - AI工具安全审计与认证平台
After=network.target

[Service]
Type=simple
User=z
WorkingDirectory=/home/z/my-project/aishield
ExecStart=/home/z/.venv/bin/gunicorn api.server_flask:app --bind 0.0.0.0:8450 --workers 4 --timeout 30 --threads 4 --access-logfile /home/z/my-project/aishield/access.log
Restart=always
RestartSec=5
Environment=PYTHONPATH=/home/z/my-project/aishield

[Install]
WantedBy=multi-user.target
